# firework/sh/config/buildkitd.service
[Unit]
Description=Buildkit daemon

[Service]
ExecStart=/usr/local/bin/buildkitd

[Install]
WantedBy=multi-user.target